# install/touchpad-enhance.service.in
[Unit]
Description=Touchpad Enhance
Documentation=https://github.com/suif4599/touchpad-enhance
After=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/touchpad-enhance
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hotkey-manager

[Install]
WantedBy=multi-user.target